# supervisord-tomcat.conf
[program:tomcat]
command=/usr/local/bin/supervisord-tomcat.sh
numprocs=1
autostart=true
autorestart=true
priority=2
